# systemd/xflux.service
[Unit]
Description=f.lux daemon

[Service]
ExecStart=/usr/bin/xflux -z 14623 -nofork

[Install]
WantedBy=default.target
